# ucxsync.service
[Unit]
Description=UCXSync - File synchronization service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/ucxsync
Environment=UCXSYNC_SERVICE_NAME=%N
ExecStart=/opt/ucxsync/ucxsync --config /etc/ucxsync/config.yaml
Restart=on-failure
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ucxsync

# Security
NoNewPrivileges=true
# Do not isolate mounts via PrivateTmp; services must observe live /ucdata mount changes.

[Install]
WantedBy=multi-user.target
